# 04/10-keyboard.rules
ACTION=="add", \
SUBSYSTEM=="input", \
ENV{ID_BUS}="usb", \
RUN+="/sbin/modprobe /root/little-penguin/04/hello.ko"
